# examples_rpi_pc/hc-mqtt.service
[Unit]
Description=Run the HC MQTT application
After=multi-user.target

[Service]
ExecStart=FOLDER/hc-mqtt

[Install]
WantedBy=multi-user.target
